# deploy/coterie-backup.service
[Unit]
Description=Coterie database backup
Documentation=man:sqlite3(1)
# Order after coterie.service if both happen to start together (boot).
# We do NOT Requires= it — VACUUM INTO works against a stopped DB
# just fine, and a backup that fails because the main service is down
# for maintenance is a worse outcome than a backup that runs anyway.
After=coterie.service

[Service]
Type=oneshot
User=coterie
Group=coterie

# The script reads env from /etc/default/coterie-backup if present.
# Operators typically set COTERIE_BACKUP_S3_URI and AWS credentials there.
EnvironmentFile=-/etc/default/coterie-backup

ExecStart=/opt/coterie/deploy/backup.sh

# The backup runs unattended; surface failures via systemd's status
# rather than swallowing them. journalctl -u coterie-backup will show
# the script's stdout/stderr.
StandardOutput=journal
StandardError=journal

# Same hardening as the main service, scoped to what the script needs.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/coterie/backups

[Install]
WantedBy=multi-user.target
